# dping@.service
[Unit]
Description=Ping camera on %I

[Service]
ExecStart=/bin/sh -c 'dping -i %I'

[Install]
WantedBy=dtools.target
