# dpoll@.service
[Unit]
Description=Get videos from camera on %I

[Service]
ExecStart=~/.local/bin/dpoll -i %I
#[Install]
#WantedBy=multi-user.target
#Also=gsl.timer