# dpoll@.service
[Unit]
Description=Get videos and images from camera on %I

[Service]
ExecStart=/bin/sh -c 'dpoll -i %I'
#[Install]
#WantedBy=dtools.target
