# dtools.target
[Unit]
Description=Start record, poll  and ping camera at user login
[Install]
WantedBy=default.target
