# gsl.timer
[Unit]
Description=Get sleep video from DCS-936L camera each hour

[Timer]
Persistent=true
OnCalendar=hourly

[Install]
WantedBy=network-online.target
